# acrn-board-parser.service
[Unit]
After=network.target

[Service]
ExecStart=/opt/workspace/acrn-board-parser.sh

[Install]
WantedBy=default.target
